# s6-rc.d/hishtory/up
/etc/s6-overlay/s6-rc.d/hishtory/scripts/run.sh