# s6-rc.d/chezmoi/up
/etc/s6-overlay/s6-rc.d/chezmoi/scripts/run.sh